# dns_updater/dns_updater.conf
[program:dns_updater]
command=python /dns_updater/gandi-live-dns.py
directory=/
user=ubuntu
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true