# dns_updater/dns_updater.conf
[program:tripwatcher]
command=python /home/gandi-live-dns.py
directory=/home
user=ubuntu
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true