# dns_updater/dns_updater.conf
[program:dns_updater]
command=python /install/dns_updater/gandi-live-dns.py
directory=/install
user=ubuntu
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true